# correosfunc/crontab/maillinux.sh
#!/bin/bash
# -*- ENCODING: UTF-8 -*
cd /opt/lampp/htdocs/correosfundacion/Mails/
/opt/lampp/bin/php /opt/lampp/htdocs/correosfundacion/Mails/mail.php
exit
